# dist/keald.local.service
[Unit]
Description=keal system daemon
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/keald
Restart=on-failure

[Install]
WantedBy=multi-user.target
